# www/crontab
0,15,30,45 * * * * /Users/c_dignazio_laptop/Desktop/catherine/projects_art/babblingbrook/TheBabblingBrook/www/saveWeatherData.sh
